# supervisord_squid.conf
[program:squid]
command = squid -N -f /etc/squid/squid.conf
autorestart = true
stderr_logfile = /tmp/squid.err
stdout_logfile = /tmp/squid.log
